# deploy/cockpit-collector.service
[Unit]
Description=Cockpit Collector Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cockpit
Group=cockpit
EnvironmentFile=/etc/cockpit/cockpit.env
ExecStart=/usr/local/bin/cockpit-collector
Restart=always
RestartSec=10
CPUQuota=50%
MemoryMax=256M
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
